# tiktok-vid2seq.service
[Unit]
Description=TikTok Vid2Seq BLIP2 Service
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/tiktok_production
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
EnvironmentFile=/home/user/tiktok_production/fix_ffmpeg_env.sh
ExecStart=/usr/bin/python3 /home/user/tiktok_production/api/pytorch_vid2seq_optimized.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target